# supervisor/suap.conf
[program:suap]
directory = /opt/suap/
user = www-data
command = /opt/scripts/run_suap.sh
stdout_logfile = /opt/logs/suap_out.log
stderr_logfile = /opt/logs/suap_err.log
redirect_stderr=true
